# systemd/shopee-warmer.service
[Unit]
Description=Shopee Session Warmer
After=network.target

[Service]
Type=simple
User=akbar
WorkingDirectory=/home/akbar/task-weekly
EnvironmentFile=/home/akbar/task-weekly/shopee-session-monitor/.env
Environment=TZ=Asia/Jakarta
Environment=HEADLESS=true
Environment=PYTHONUNBUFFERED=1
Environment=OFD_DISCORD_MODE=1

# Jalankan warmer langsung dengan uv
ExecStart=/usr/local/bin/uv run --project src shopee-session-monitor/warmer.py

# Batasi memory — Chrome bisa makan 400-600MB
MemoryMax=800M
MemoryHigh=700M

# Restart otomatis jika crash
Restart=always
RestartSec=15

# Log output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shopee-warmer

[Install]
WantedBy=multi-user.target
